# k8s-queue-monitor/k8s-monitor.crontab
*/15 * * * * /home/USERNAME/k8s-queue-monitor/wrapper.sh
