# k8s-queue-monitor/k8s-monitor.crontab
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

K8S_QUEUE_MONITOR_KUBECTL_PATH=/usr/local/bin/kubectl
K8S_QUEUE_MONITOR_KUBECONFIG=/home/USERNAME/.kube/config
K8S_QUEUE_MONITOR_OUTPUT_DIR=/home/USERNAME/k8s-queue-monitor-data

*/15 * * * * /home/USERNAME/k8s-queue-monitor/wrapper.sh
